# roles/memoria-core/files/etc/cron.d/docker
0 7 * * *           root    docker restart $(docker ps -q)
